# target/linux/rockchip/armv8/base-files/etc/init.d/firefly
#!/bin/sh /etc/rc.common
START=19

start()
{                                                                                                                                                                     
        [ -f /etc/config/wireless ] && return 0
        wifi config
}
